# osquery/2.Windows/windows_environment_check.conf
{
  "platform":"windows",
  "description":"ATT&CK: T1053",
  "queries":{
    "Scheduled_Task":{
      "query":"select name,action,path,enabled,datetime(next_run_time,'unixepoch','UTC') from scheduled_tasks;",
      "interval":600,
      "description":"Lists all of the tasks in the Windows task scheduler - ATT&CK T1053",
      "platform":"windows",
      "removed":false
    },
    "Snaphost_Scheduled_Task":{
      "query":"select name,action,path,enabled,datetime(next_run_time,'unixepoch','UTC') from scheduled_tasks;",
      "interval":3600,
      "description":"Lists all of the tasks in the Windows task scheduler",
      "platform":"windows",
      "snapshot":true
    },
    "Services":{
      "query":"SELECT * FROM services WHERE start_type='DEMAND_START' OR start_type='AUTO_START';",
      "interval":600,
      "description":"Lists all installed services configured to start automatically at boot - ATT&CK T1050",
      "removed":false
    },
    "Snapshot_services":{
      "query":"SELECT * FROM services;",
      "interval":3600,
      "description":"Snapshot Services query",
      "snapshot":true
    },
    "Startup_Items":{
      "query":"select * from startup_items;",
      "interval":900,
      "description":"Windows Startup Items Osquery predefined query - ATT&CK T1060",
      "removed":false
    },
    "Snaphost_Startup_Items":{
      "query":"select * from startup_items;",
      "interval":3600,
      "description":"Snapshot Windows Startup Items Osquery predefined query - ATT&CK T1060",
      "snapshot":true
    }
  }
}
